# systemd/power.timer
[Unit]
Description=Run once a minute

[Timer]
OnUnitActiveSec=60s
OnBootSec=20s

[Install]
WantedBy=timers.target
